# cryptotech/debian/cryptotech.cron.d
0 * * * *	root	/bin/bash /hive/bin/cryptotech/ping.sh > /var/log/cryptotech-ping.log 2>&1
